# docker/backup/crontab
# Daily backup at 2 AM UTC
0 2 * * * /usr/local/bin/backup-database.sh >> /var/log/backup.log 2>&1
